# mysql/src/s6/update-daemon/run
#!/bin/sh

exec s6-setuidgid nobody php /var/www/tt-rss/update_daemon2.php
